# concourse/install/concourse-web.service
[Unit]
Description=concourse-web
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/concourse-web

[Install]
WantedBy=multi-user.target
