# concourse/concourse-web.service
[Unit]
Description=concourse-web

[Service]
ExecStart=/usr/bin/concourse-web

[Install]
WantedBy=network-online.target
